# templates/etc/cron.d/drupal.j2
# {{ ansible_managed }}

0  *   * * *   {{ drupal_user }}    curl {{ drupal_domain }}/cron.php?cron_key={{ drupal_cron_key }}  >> /dev/null 2>&1

