# templates/etc/cron.d/drupal.j2
# {{ ansible_managed }}

0  *   * * *   {{ drupal_user }}    /usr/bin/env PATH={{ drupal_cron_path }} COLUMNS=72 drush --root={{ drupal_install_path }} --uri={{ drupal_domain }} --quiet cron

